# djangotest/globaldomain/gd/xx.cron
* */1 * * * python /djangotest/globaldomain/gd/dnspython.py

